# rootfs/etc/services.d/landingpage/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Home Assistant landingpage service
# ==============================================================================

exec landingpage 
